# tests/testenv/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:krb5kdc]
command=/usr/sbin/krb5kdc -n

[program:kadmin]
command=/usr/sbin/kadmind -nofork

[program:slapd]
command=/bin/bash -c "export KRB5_KTNAME='FILE:/etc/openldap/ldap.keytab' && exec /usr/sbin/slapd -h 'ldap:// ldaps://' -d 1"

#[program:bind9]
#command=/usr/sbin/named -c /etc/named.conf -f -u named
